# debian/xscript-standard.files
usr/lib/xscript/xscript*so*
usr/bin/xscriptcache*
etc/cron.d/xscript-cache*
